# templates/default/supervisor_netconfig.conf.erb
[program:netconfig]
command = gunicorn app:app -b localhost:<%= @port %>
directory = <%= @home %>
user = <%= @user %>
